# roles/backup-server/templates/rsnapshot-hourly.timer.j2
[Unit]
Description=Rsnapshot hourly timer: {{ item.name }}

[Timer]
OnCalendar=hourly
RandomizedDelaySec=2700

[Install]
WantedBy=multi-user.target
