# roles/backup-server/templates/rsnapshot-daily.timer.j2
[Unit]
Description=Rsnapshot daily timer: {{ item.name }}

[Timer]
OnCalendar=*-*-* 00:00:00
RandomizedDelaySec=3600

[Install]
WantedBy=multi-user.target
