# roles/aide/files/aide-check.service
[Unit]
Description=AIDE Check

[Service]
Type=oneshot
ExecStart=/usr/sbin/aide --check